# etc/systemd/user/autostart-exam.service
[Service]
ExecStart=/usr/bin/lernstick-rdp-start-exam
