# etc/systemd/user/autostart-exam.path
[Unit]
Description="Exam autostart"

[Path]
PathExists=/home/user/Schreibtisch/exam.zip

[Install]
WantedBy=default.target
